# com.txf.inventory.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.txf.inventory</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/coleharing/.nvm/versions/node/v22.14.0/bin/node</string>
        <string>/Users/coleharing/Cursor_Projects/TXF_Automations/src/automation/daily-inventory-run.js</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/coleharing/Cursor_Projects/TXF_Automations</string>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/coleharing/Cursor_Projects/TXF_Automations/logs/daily-run.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/coleharing/Cursor_Projects/TXF_Automations/logs/daily-run-error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/coleharing/.nvm/versions/node/v22.14.0/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>